# deploy/actio-brand-budget-weekday.timer
[Unit]
Description=Trigger BRAND weekday budget (30 zł) — every Monday 06:00

[Timer]
OnCalendar=Mon *-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
