# deploy/actio-brand-budget-weekday.service
[Unit]
Description=Actio BRAND budget — weekday (30 zł)
After=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/actio-marketing-panel
ExecStart=/root/.local/bin/uv run python scripts/set_brand_budget.py 30
StandardOutput=append:/opt/actio-marketing-panel/logs/brand-budget.log
StandardError=append:/opt/actio-marketing-panel/logs/brand-budget.err.log
