# deploy/actio-autopublish-checker.timer
[Unit]
Description=Run Actio Autopost IMAP checker every 30 min

[Timer]
OnCalendar=*:0/30
Persistent=true

[Install]
WantedBy=timers.target
