# deploy/actio-brand-budget-weekend.timer
[Unit]
Description=Trigger BRAND weekend budget (60 zł) — every Friday 06:00

[Timer]
OnCalendar=Fri *-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
