# deploy/actio-autopublish-generate.timer
[Unit]
Description=Run Actio Autopost generator (wt + czw 9:00)

[Timer]
OnCalendar=Tue,Thu 09:00
Persistent=true

[Install]
WantedBy=timers.target
