# deploy/actio-daily-report.timer
[Unit]
Description=Actio Marketing — Daily Report at 7:00 (Europe/Warsaw)

[Timer]
OnCalendar=*-*-* 07:00:00
Persistent=true

[Install]
WantedBy=timers.target
