# deploy/sendly-geo-monitor.timer
[Unit]
Description=Sendly — GEO monitor (co tydzien, poniedzialek)

[Timer]
OnCalendar=Mon *-*-* 06:30 Europe/Warsaw
Persistent=true

[Install]
WantedBy=timers.target
