# deploy/sendly-daily-report.timer
[Unit]
Description=Sendly — Daily Report at 07:01 (Europe/Warsaw)

[Timer]
# 07:01 — tuz po actio (07:00), zeby raporty szly po sobie. Jawna strefa: niezalezne od TZ hosta.
OnCalendar=*-*-* 07:01 Europe/Warsaw
Persistent=true

[Install]
WantedBy=timers.target
